# deploy/supervisord.conf
[supervisord]
nodaemon=true

[program:bun]
command=bun /app/openclaude/scripts/start-http-stream.ts
environment=HTTP_PORT=50052,HOST=127.0.0.1
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nest]
command=node /app/back/dist/src/main.js
directory=/app/back
autorestart=true
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0